# supervisord.conf
[supervisord]
nodaemon=true

[program:Rserve]
command=/opt/R/bin/R CMD Rserve --quiet --vanilla --RS-conf /etc/Rserve.conf

[program:postgres]
command=/usr/lib/postgresql/9.3/bin/postgres -D /var/lib/postgresql/9.3/main -c config_file=/etc/postgresql/9.3/main/postgresql.conf
user=postgres

[program:tomcat7]
command=/usr/share/tomcat7/bin/catalina.sh run
environment=CATALINA_BASE=/var/lib/tomcat7,CATALINA_TMPDIR=/var/lib/tomcat7/tmp
user=tomcat7
group=tomcat7

[program:solr]
command=java -jar start.jar
directory=/home/transmart/transmart-data/solr
user=transmart



